# aiyawocao/aiyawocao-web/src/main/service/skrbt.service
[Unit]
Description=skrbt
[Service]
ExecStart=start_bash_path
ExecStop=stop_bash_path
Restart=always
RestartSec=5
[Install]
WantedBy=default.target
